# trading-linux/buildroot-external/board/trading/overlay/etc/systemd/system/multi-user.target.wants/trading-system.service
[Unit]
Description=FPGA Trading System Orchestrator
After=network.target xdma.service nvidia-gpu.service
Requires=xdma.service
Wants=nvidia-gpu.service

[Service]
Type=simple
ExecStart=/opt/trading/bin/trading_system_orchestrator /opt/trading/config/system_config.json
ExecStop=/bin/kill -SIGTERM $MAINPID
WorkingDirectory=/opt/trading
Restart=on-failure
RestartSec=5
User=root

# Performance tuning
Nice=-20
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=90
LimitMEMLOCK=infinity
LimitRTPRIO=99

# Environment
Environment="LD_LIBRARY_PATH=/opt/cuda/lib64:/opt/xgboost/lib"
Environment="CUDA_VISIBLE_DEVICES=0"

[Install]
WantedBy=multi-user.target
