# trading-linux/buildroot-external/board/trading/overlay/etc/systemd/system/nvidia-cpu-pin.service
[Unit]
Description=Pin NVIDIA processes to system CPUs
After=nvidia-modprobe.service
Wants=nvidia-modprobe.service

[Service]
Type=oneshot
ExecStart=/opt/trading/bin/pin_nvidia.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

